# power_gadget.socket
[Unit]
Description=Prometheus metric for cpu wattage

[Socket]
ListenStream=9999
Accept=yes

[Install]
WantedBy=sockets.target
