# power_gadget@.service
[Unit]
Description=Prometheus metric for cpu wattage

[Service]
ExecStart=-/usr/sbin/power_gadget
StandardOutput=socket